# roles/docker_compose/templates/osism.target.j2
[Unit]
Description=OSISM target allowing to start/stop all OSISM service at once

[Install]
WantedBy=default.target
